# debian/postrm
#!/bin/sh
# postrm script for chrony
#
# see: dh_installdeb(1)

set -e

# targets: purge|remove|upgrade|failed-upgrade|abort-install|abort-upgrade|disappear

restore_timesyncd() {
    # on next reboot it would start, but that would leave time
    # unsynchronized until then. So as the Conflicts in the service file kill
    # systemd-timesyncd re-establish it if it is enabled
    if [ "$(systemctl is-enabled systemd-timesyncd 2>/dev/null)" = "enabled" ] ; then
        systemctl start systemd-timesyncd
    fi
}

case "$1" in
    purge)
        rm -f /var/lib/chrony/*
        rm -f /etc/chrony/chrony.conf
        rm -f /etc/chrony/chrony.keys
        if `which ucf >/dev/null`
        then
            ucf --purge /etc/chrony/chrony.conf
            ucf --purge /etc/chrony/chrony.keys
            if [ -x "$(which ucfr)" ]; then
                ucfr --purge chrony /etc/chrony/chrony.conf
                ucfr --purge chrony /etc/chrony/chrony.keys
            fi
        fi
        rm -rf /etc/chrony
        rm -rf /run/chrony || true
        rm -rf /var/lib/chrony
        rm -rf /var/log/chrony
        # Remove "_chrony" system user/group
        if [ -x "$(which deluser)" ]
        then
            deluser --quiet --system _chrony > /dev/null 2>&1 || true
        fi

        restore_timesyncd
    ;;

    remove)
        restore_timesyncd
    ;;

    upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)

    ;;

    *)
        echo "postrm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
